# spotify/rootfs/etc/services.d/spotifyd/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Spotify Connect
# Runs librespot
# ==============================================================================
declare -a options
declare bitrate
declare name
declare oauth
declare autoplay
declare initial_volume

bashio::log.info 'Starting the Spotify daemon...'

# Bitrate
options+=(--bitrate $(bashio::config 'bitrate'))

# Initial Volume
initial_volume=$(bashio::config 'initial_volume')
options+=(--initial-volume "$initial_volume")

# Device name
name=$(bashio::config 'name')
options+=(--name "${name}")

# Autoplay
if bashio::config.true 'autoplay'; then
  options+=(--autoplay on)
fi

# oAuth
if bashio::config.true 'oauth'; then
  options+=(--enable-oauth)
fi

# Save writes
options+=(--disable-audio-cache)
options+=(--system-cache /config)

# Are we running in debug mode?
if bashio::debug; then
  options+=(--verbose)
fi

# Run librespot
exec librespot "${options[@]}"
